# .config/systemd/user/bpd.service
[Unit]
Description=Beets Music Player Daemon

[Service]
ExecStart=/usr/bin/beet bpd

[Install]
WantedBy=default.target