# .config/systemd/user/pacman-user-backup.service
[Unit]
Description=Pacman and Cower package list generation scripts

[Service]
Type=oneshot
ExecStart=/home/stoye/dotfiles/update-lists.sh